# claude-rc.service
[Unit]
Description=Claude RC Launcher
After=network.target

[Service]
Type=simple
EnvironmentFile=%h/.claude-rc/env
ExecStart=/usr/bin/env python3 %h/.claude-rc/app/app.py
Restart=on-failure
RestartSec=5
StandardOutput=append:%h/.claude-rc/logs/claude-rc.log
StandardError=append:%h/.claude-rc/logs/claude-rc.log

[Install]
WantedBy=default.target
