# claude-rc.service
[Unit]
Description=Claude RC Launcher
After=network.target

[Service]
Type=simple
# Only kill the launcher itself on stop/restart — the tmux server (and every
# Claude session in it) lives in the same cgroup and must survive updates.
KillMode=process
EnvironmentFile=%h/.claude-rc/env
ExecStart=/usr/bin/env python3 %h/.claude-rc/app/app.py
Restart=on-failure
RestartSec=5
StandardOutput=append:%h/.claude-rc/logs/claude-rc.log
StandardError=append:%h/.claude-rc/logs/claude-rc.log

[Install]
WantedBy=default.target
